# headsetk.service
[Unit]
After=multi-user.target

[Service]
ExecStart=/home/tomasz/headsetk/target/release/headsetk

[Install]
WantedBy=graphical.target
